# server/vpn_bot.service
[Unit]
Description=VPN Telegram Bot
After=network.target

[Service]
User=root
WorkingDirectory=/root
ExecStart=/usr/bin/python3 /root/vpn_bot.py
Restart=always
RestartSec=10s
MemoryMax=128M

# Токен и admin ID — НЕ в коде, а здесь
# Перевыпуа сти токен у BotFather и вставь сюда
Environment=BOT_TOKEN=REPLACE_WITH_YOUR_BOT_TOKEN
Environment=ADMIN_ID=REPLACE_WITH_YOUR_ADMIN_ID

[Install]
WantedBy=multi-user.target
