# server/vpn_bot.service
[Unit]
Description=VPN Telegram Bot
After=network.target

[Service]
User=root
WorkingDirectory=/root
ExecStart=/usr/bin/python3 /root/vpn_bot.py
Restart=always
RestartSec=10s
MemoryMax=128M

# Токен и admin ID — НЕ в коде, а здесь
Environment=BOT_TOKEN=8675996023:AAG23QhfNxTbMzVwYUeSIc6AMceJz9urXeQ
Environment=ADMIN_ID=REPLACE_WITH_YOUR_ADMIN_ID

[Install]
WantedBy=multi-user.target
